# 03-4-ansible/files/3_/roles/wcg/templates/wcg.service.j2
[Unit]
Description=wcg service

[Service]
ExecStart={{ wcg_path }}/word-cloud-generator
Restart=always

[Install]
WantedBy=multi-user.target
